# debian13/overlay-sb-debian/home/kiosk/.config/systemd/user/kiosk-session.target
[Unit]
Description=Kiosk session target
BindsTo=kiosk-session.target
Before=kiosk-session.target

