# debian13/misc/overlay/etc/xdg/autostart/kiosk.desktop
[Desktop Entry]
Type=Application
Name=kiosk
Exec=gnome-terminal
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
