# debian13/overlay-installer/root/.config/systemd/user/weston.socket
[Unit]
Description=Weston socket

[Socket]
ListenStream=%t/wayland-0

